# config/supervisor.conf
[program:md]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.py
directory=/root/ResumeCenter
autostart=true
autorestart=true
